# supervisor/conf.d/processes.conf.ini
[program:testing_app]
command                 = /usr/bin/python3 /testing_app/%(program_name)s.py

directory               = /testing_app/
user                    = root

autostart               = true
startsecs               = 10
startretries            = 100
autorestart             = true

stdout_logfile          = /testing_app/logs/%(program_name)s.stdout
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups  = 10

stderr_logfile          = /testing_app/logs/%(program_name)s.stderr
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups  = 10
